# setup/99-usbftdi.rules
ATTRS{manufacturer}=="BertholdTech", RUN+="/etc/mountBerthold.sh"
